# qdo.conf
description "QDo daemon"

start on filesystem and net-device-up IFACE!=lo
stop on runlevel [!2345]

respawn

script
    /usr/bin/qdo -f /var/qdo/ -r /opt/qdo/
end script
